# postfix/supervisord.conf
[supervisord]
nodaemon=true

[program:rsyslog]
command=/usr/sbin/rsyslogd -n

[program:postfix]
command=/usr/lib/postfix/master -d -c /etc/postfix

[program:saslauth]
command=/usr/sbin/saslauthd -a shadow -c -d -m /var/spool/postfix/var/run/saslauthd -n 5
