# cron/diff-backup.cron
0 4 * * * @backup@/diff-backup.sh source=@source@ dest=@dest@ 1>&2
