# cron/docker-backup.cron
0 2 * * * @backup@/docker-backup.sh local=@local@/ 1>&2
