# cron/remote-backup.cron
0 3 * * * /usr/bin/rsync -aP --delete @source@/ @dest@/ 1>&2
