# bind/supervisord.conf
[supervisord]
nodaemon=true

[program:rsyslog]
command=/usr/sbin/rsyslogd -n

[program:bind]
command=/usr/sbin/named -f -u bind
